# client/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server --dir /app/data/redis --save 60 1 --appendonly yes --appendfsync everysec
user=racksmith
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app]
command=uv run uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app
user=racksmith
environment=HOME="/app"
autorestart=true
startsecs=5
startretries=5
stopsignal=INT
stopwaitsecs=30
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=uv run python run_worker.py worker.settings.WorkerSettings
directory=/app
user=racksmith
environment=HOME="/app"
autorestart=true
startsecs=5
startretries=5
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
user=racksmith
autorestart=true
stopsignal=QUIT
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
